# configs/provision/client.d/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
chmod=0700                 ; socket file mode (default 0700)
chown=vagrant       ; socket file uid:gid owner

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
user=vagrant                 ; (default is current user, required if root)
environment=DBUS_SESSION_BUS_ADDRESS="unix:path=/run/dbus/system_bus_socket"     ; (key value pairs to add to environment)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

#[program:jackd]
#;command=/usr/bin/jackd -dalsa -dhw:USB -r48000 -p1024 -n2
#;command=/usr/bin/jackd -ddummy -r48000 -p1024    ; the program (relative uses PATH, can take args)
#command=/home/vagrant/start_jackd.sh
#priority=599                  ; the relative start priority (default 999)
#stopasgroup=true

[program:client]
;command=/usr/local/bin/matriz -d
command=python /home/vagrant/territorios/src/client.py -d --alsa
directory=/home/vagrant
environment=HOME="/home/vagrant"
priority=699                  ; the relative start priority (default 999)
autostart=true                ; start at supervisord start (default: true)
autorestart=true
stopasgroup=true             ; send stop signal to the UNIX process group (default false)

[program:shutdown_pi]
command=python /home/pi/shutdown_pi.py
priority=699                  ; the relative start priority (default 999)
autostart=true                ; start at supervisord start (default: true)
autorestart=true
