# configs/supervisor/conf.d/matriz.conf
[program:matriz]
command=//home/pi/territorios/src/server.py
priority=550
autostart=false
user=pi
